# root/etc/crontabs/filebot
* * * * * /scripts/filebot_import.sh
0 3,15 * * * /scripts/filebot_artwork.sh
